# launchd/com.codex.goal-status-hook.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.codex.goal-status-hook</string>

  <key>ProgramArguments</key>
  <array>
    <string>python3</string>
    <string>__HOME__/.codex/goal-status-hook/goal_status_report.py</string>
    <string>--watch</string>
    <string>--interval</string>
    <string>0.25</string>
  </array>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__HOME__/.codex/goal-status-hook/launchd.out.log</string>

  <key>StandardErrorPath</key>
  <string>__HOME__/.codex/goal-status-hook/launchd.err.log</string>
</dict>
</plist>
